# cookbooks/crontab/files/default/crontab
0 0 * * * root /usr/local/bin/rvm-shell '1.9.2' -c 'gem mirror'
