# cookbooks/crontab/files/default/crontab
0 0 * * * root /usr/local/bin/rvm-shell '1.9.2' -c 'gem mirror' 2>&1 >/dev/null
0 10 * * * root /usr/local/bin/rvm-shell '1.9.2' -c '/root/gem-sniffer' 2>&1 >/dev/null
